# little/task_dir_05/task05.rules
SUBSYSTEM=="input", SYMLINK+="task05", RUN+="/usr/bin/install.sh"
